# build/package/newrelic-super-agent.service
[Unit]
Description=New Relic Super Agent
After=network.target

[Service]
EnvironmentFile=/etc/newrelic-super-agent/newrelic-super-agent.conf
ExecStart=/usr/bin/newrelic-super-agent
KillMode=mixed
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
