# build/package/newrelic-agent-control.service
[Unit]
Description=New Relic Agent Control
After=network.target

[Service]
# kept for backward compatibility. The environment_variables.yaml is used instead for new installs.
EnvironmentFile=/etc/newrelic-agent-control/systemd-env.conf
ExecStart=/usr/bin/newrelic-agent-control
KillMode=mixed
Restart=on-failure
Type=simple

[Install]
WantedBy=multi-user.target
